# systemd/go-aws-mon.service
[Unit]
Description=go-aws-mon

[Service]
ExecStart=/home/core/go-aws-mon.sh

[Install]
WantedBy=multi-user.target
